# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling WAN tunnel."
    systemctl enable waggle-wan-tunnel.service
    ;;
esac
