# deb/install/prerm
#!/bin/bash

systemctl stop waggle-wan-tunnel.service || true
systemctl disable waggle-wan-tunnel.service || true
