# ROOTFS/etc/systemd/system/waggle-wan-tunnel.service
[Unit]
Description=Waggle WAN Tunnel

[Service]
Type=notify
Restart=always
RestartSec=10s
ExecStart=/usr/bin/waggle_wan_tunnel.py

[Install]
WantedBy=multi-user.target
